# com.filemodtracker.ui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.filemodtracker.ui</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/filemodtracker</string>
        <string>ui</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/filemodtracker-ui.err</string>
    <key>StandardOutPath</key>
    <string>~/Library/Logs/filemodtracker-ui.log</string>
</dict>
</plist>