# com.filemodtracker.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.filemodtracker.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/filemodtracker</string>
        <string>daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/var/log/filemodtracker.err</string>
    <key>StandardOutPath</key>
    <string>/var/log/filemodtracker.log</string>
</dict>
</plist>